# configs/conf.d/redis.conf
[program:redis]
command=/usr/local/bin/redis-server

user=root
numprocs=1
stdout_logfile=/var/log/redis/redis.log
stderr_logfile=/var/log/redis/redis_err.log
autostart=true
autorestart=true
startsecs=10

stopasgroup=true
priority=998
